# ebuilds/ebuilds/ebuilds/net-mail/messagewall/files/messagewall
#!/sbin/runscript
# Copyright 2003 Alexander Holler
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /home/xubuntu/berlios_backup/github/tmp-cvs/gentoo-deutsch/Repository/ebuilds/ebuilds/ebuilds/net-mail/messagewall/files/Attic/messagewall,v 1.1 2003/08/22 15:27:29 fow0ryl Exp $

opts="${opts} stats"

depend() {
	need net
	after postfix sendmail qmail exim
	use dns logger
}

checkconfig() {
	if [ ! -f /etc/messagewall/messagewall.conf ] ; then
		eerror "You need to setup /etc/messagewall/messagewall.conf first!"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting MessageWall"
	/usr/bin/messagewall 2>&1 | /usr/bin/logger -t MessageWall -p mail.info &
	eend $?
}

stop() {
	ebegin "Stopping MessageWall"
	/usr/bin/messagewallctl stop
	eend $?
}

stats() {
	/usr/bin/messagewallstats /var/log/mail/messagewall.log
}
