# ebuilds/net-dns/rbldnsd/files/rc
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /home/xubuntu/berlios_backup/github/tmp-cvs/gentoo-deutsch/Repository/ebuilds/net-dns/rbldnsd/files/rc,v 1.2 2004/01/31 14:59:19 mad Exp $

depend() {
	need net 
}

start() {
	ebegin "Starting rbldnsd-server"
	start-stop-daemon --start --quiet --pidfile /var/run/rbldnsd.pid --exec /usr/sbin/rbldnsd -- ${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping rbldnsd-server"
	start-stop-daemon --stop --quiet --pidfile /var/run/rbldnsd.pid --exec /usr/sbin/rbldnsd
	eend $?
}
